# systemd/av-scan.timer
# Rendered by `av-scan --install-timer` (placeholders @BOOTDELAY@/@INTERVAL@/@NAME@).
#   OnBootSec       = @BOOTDELAY@ after every boot -> catches machines that are
#                     powered off unpredictably (a scan runs on each startup).
#   OnUnitActiveSec = @INTERVAL@, derived from the millisecond value you passed
#                     (--install-timer 86400000 -> 86400s = 1 day) -> repeats
#                     that often while the machine stays on.
# Persistent=true also re-runs a run that came due while the box was off.
[Unit]
Description=linux-av scheduled scan timer
Documentation=https://github.com/danganhtu01/linux-av

[Timer]
OnBootSec=@BOOTDELAY@
OnUnitActiveSec=@INTERVAL@
AccuracySec=1s
Persistent=true
Unit=@NAME@.service

[Install]
WantedBy=timers.target
